# obexd/src/obex.service.in
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=/bin/sh -c 'if [ -z `ps ax | grep -v grep | grep obexd` ] ; then exec @libexecdir@/obexd -d --noplugin=syncevolution,pcsuite,irmc --symlinks -r /opt/usr/media; fi'

[Install]
Alias=dbus-org.bluez.obex.service
