# obexd/src/obex.service.in
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=/bin/obex.sh

[Install]
Alias=dbus-org.bluez.obex.service
